# site.conf
# this goes in /etc/init/
# to start run 
# sudo initctl start site

start on startup
chdir /home/mark/stocks/littlejohn
exec /home/mark/.cabal/bin/site -p 80
